# laraver_queue.ini
[program:cec_queue]
command                 = /usr/bin/php  /data/www/Cec/current/artisan queue:work --sleep=3 --tries=3 --daemon
directory               = /data/www/Cec/current
process_name            = %(program_name)s_%(process_num)s
user                    = root
numprocs                = 4
autostart               = true
autorestart             = true
stdout_logfile          = /data/logs/access_supervisor_Queue.log
stdout_logfile_maxbytes = 10MB
stderr_logfile          = /data/logs/error_supervisor_Queue.log
stderr_logfile_maxbytes = 10MB